# deploy/pulses/ai.flyn.pulse.conv-summarize-backfill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>           <string>ai.flyn.pulse.conv-summarize-backfill</string>
    <key>ProgramArguments</key>
    <array>
      <string>{{HOME}}/AI/openclaw/flyn-agent/deploy/pulses/conv_summarize_backfill.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key> <string>{{HOME}}</string>
      <key>PATH</key> <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <!-- Daily at 04:15. Catches summaries that fell through during the day. -->
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>   <integer>4</integer>
      <key>Minute</key> <integer>15</integer>
    </dict>
    <key>StandardOutPath</key> <string>/tmp/flyn-conv-summarize-backfill.log</string>
    <key>StandardErrorPath</key><string>/tmp/flyn-conv-summarize-backfill.log</string>
  </dict>
</plist>
