# deploy/launchd/ai.flyn.graphiti-api.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.flyn.graphiti-api</string>

  <key>ProgramArguments</key>
  <array>
    <string>{{HOME}}/.openclaw/workspace/memory/structured/graphiti-venv/bin/python</string>
    <string>{{HOME}}/.openclaw/workspace/kg/flyn-graphiti-api.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>{{HOME}}/.openclaw/workspace/kg</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>{{HOME}}/.openclaw/workspace/memory/structured/graphiti-venv/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>{{HOME}}</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/tmp/flyn-graphiti-api.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/flyn-graphiti-api.err</string>
</dict>
</plist>
