# deploy/launchd/ai.flyn.pulse.meeting-categorize.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>ai.flyn.pulse.meeting-categorize</string>
  <key>ProgramArguments</key><array>
    <string>/bin/bash</string>
    <string>/Users/4c/AI/openclaw/flyn-agent/deploy/cron/scripts/meeting-categorize.sh</string>
  </array>
  <key>EnvironmentVariables</key><dict>
    <key>HOME</key><string>/Users/4c</string>
    <key>PATH</key><string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>StartCalendarInterval</key><dict>
    <key>Hour</key><integer>2</integer>
    <key>Minute</key><integer>30</integer>
  </dict>
  <key>StandardOutPath</key><string>/Users/4c/.openclaw/logs/cron-ai.flyn.pulse.meeting-categorize.log</string>
  <key>StandardErrorPath</key><string>/Users/4c/.openclaw/logs/cron-ai.flyn.pulse.meeting-categorize.err</string>
</dict></plist>
